# buildroot/overlay/etc/profile.d/garage.sh
export GARAGE_WEBHOOK="https://example.com/"
